# files/hooks/container-getty.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear --skip-login --keep-baud 115200,38400,9600 pts/%I $TERM
