# files/hooks/getty.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear --skip-login %I $TERM
TTYVTDisallocate=no
