# files/hooks/gpm.service
[Unit]
Description=Console mouse daemon
ConditionPathExists=/dev/input/mice

[Service]
Type=forking
ExecStart=/usr/sbin/gpm -m /dev/input/mice -t exps2
PIDFile=/var/run/gpm.pid
