# files/hooks/gpm.service
[Unit]
ConditionPathExists=/dev/input/mice
