# files/hooks/ssh.service
[Service]
ExecStartPre=
ExecStartPre=/usr/lib/{{ PRODUCT_ID }}/sshd-keygen
ExecStartPre=/usr/sbin/sshd -t
