# files/hooks/udisks2.service
[Install]
WantedBy=
WantedBy=multi-user.target
